# vps/sureedge-api.service
[Unit]
Description=SureEdge Market API (FastAPI + Uvicorn)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/sureedge
ExecStart=/opt/sureedge/venv/bin/uvicorn api:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
